# steam-screenshot-sync.path
[Unit]
Description=Watch Steam screenshot directories for changes (Legacy version - DISABLED)
Documentation=https://github.com/PietHugen/SteamDeckScreenshotSync

[Path]
# Monitor the screenshots datafile for changes. The longer number is probably variable so could be different between steamdecks
PathChanged=%h/.local/share/Steam/userdata/865713241/760/screenshots.vdf

[Install]
WantedBy=default.target
