# steam-screenshot-sync.path
[Unit]
Description=Watch Steam screenshot directories for changes
Documentation=https://github.com/PietHugen/SteamDeckScreenshotSync

[Path]
# Monitor Steam's screenshot directories specifically
PathModified=%h/.local/share/Steam/userdata/*/760/remote/*/screenshots

[Install]
WantedBy=default.target
