# steam-screenshot-sync.service
[Unit]
Description=Steam Screenshot Auto-Sync with Game Names
Documentation=https://github.com/yourusername/steam-screenshot-sync
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/bin/steam-screenshot-sync.sh

# Environment variables (if needed)
Environment=PATH=%h/bin:/usr/local/bin:/usr/bin:/bin

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/rclone %h/.cache/rclone

[Install]
WantedBy=default.target